# .build/DEB/postinst
#!/bin/sh
set -e

# Source debconf library
. /usr/share/debconf/confmodule

# Create rpxy user if it doesn't exist
if ! getent passwd rpxy > /dev/null; then
    adduser --system --group --no-create-home --shell /usr/sbin/nologin rpxy
fi

# Set correct ownership for config directory
if [ -d /etc/rpxy ]; then
    chown -R rpxy:rpxy /etc/rpxy
fi

# Reload systemd, enable and start the service
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
    deb-systemd-helper unmask rpxy.service >/dev/null || true
    if deb-systemd-helper --quiet was-enabled rpxy.service; then
        deb-systemd-helper enable rpxy.service >/dev/null || true
    else
        deb-systemd-helper update-state rpxy.service >/dev/null || true
    fi
    if [ -d /run/systemd/system ]; then
        systemctl --system daemon-reload >/dev/null || true
        if [ -n "$2" ]; then
            deb-systemd-invoke try-restart rpxy.service >/dev/null || true
        else
            deb-systemd-invoke start rpxy.service >/dev/null || true
        fi
    fi
fi

exit 0
