# .build/DEB/prerm
#!/bin/bash

set -e

# Stop and disable the service before removing
systemctl stop rpxy || true
systemctl disable rpxy || true
systemctl daemon-reload

exit 0
