# .build/DEB/prerm
#!/bin/bash

set -e

systemctl stop rpxy || true
systemctl disable rpxy || true

exit 0